# i3/.config/i3/amdgputuner/amdgputuner.service
[Unit]
Description=Apply GPU settings

[Service]
ExecStart=/usr/bin/amdgputuner
