# i3/.config/i3/arch_pkg/amdgputuner.service
[Unit]
Description=Apply GPU settings on boot

[Service]
ExecStart=/usr/bin/amdgputuner

[Install]
WantedBy=multi-user.target
